# templates/default/sparkhistoryserver.service.erb
[Unit]
Description = Spark historyserver
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=80
StartLimitBurst=10

[Service]
User = <%= node['hops']['hdfs']['user'] %>
PIDFile = /tmp/spark-<%= node['hops']['hdfs']['user'] %>-org.apache.spark.deploy.history.HistoryServer-1.pid
ExecStart = <%= node['hadoop_spark']['base_dir'] %>/sbin/start-history-server.sh
ExecStop = <%= node['hadoop_spark']['base_dir'] %>/sbin/stop-history-server.sh
LimitNOFILE = 50000
Type = forking
Restart=on-failure
RestartSec=5

[Install]
WantedBy = multi-user.target
