# templates/default/sparkhistoryserver.service.erb
[Unit]
Description = Spark historyserver
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.hops.hdfs.user %>
PIDFile = /tmp/spark-<%= node.hops.hdfs.user %>-org.apache.spark.deploy.history.HistoryServer-1.pid
ExecStart = <%= node.hadoop_spark.home %>/sbin/start-history-server.sh
ExecStop = <%= node.hadoop_spark.home %>/sbin/stop-history-server.sh
LimitNOFILE = 30000

Type = forking

[Install]
WantedBy = multi-user.target
