# templates/default/sparkhistoryserver.service.erb
[Unit]
Description = Spark historyserver
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=30
StartLimitBurst=20

[Service]
User = <%= node['hops']['hdfs']['user'] %>
PIDFile = /tmp/spark-<%= node['hops']['hdfs']['user'] %>-org.apache.spark.deploy.history.HistoryServer-1.pid
ExecStartPre=/bin/bash -c 'while ! /usr/bin/host <%= @nn_rpc_endpoint %> ; do if systemctl list-units --full -all | grep -Fq 'namenode'; then if systemctl is-active --quiet namenode; then /bin/sleep 5; else ! break; fi; else /bin/sleep 5; fi;  done'
ExecStart = <%= node['hadoop_spark']['base_dir'] %>/sbin/start-history-server.sh
ExecStop = <%= node['hadoop_spark']['base_dir'] %>/sbin/stop-history-server.sh
LimitNOFILE = 50000
Type = forking
Restart=on-failure

[Install]
WantedBy = multi-user.target
